# data/knowledge_base/mitre_techniques.json
[
  {
    "technique_id": "T1555",
    "technique_name": "Credentials from Password Stores",
    "family": "AgentTesla",
    "text": "AgentTesla harvests stored credentials from web browsers, email clients, and FTP applications. It targets password stores in Chrome, Firefox, Outlook, and Thunderbird by reading local database files and registry entries. The stolen credentials are exfiltrated to attacker-controlled SMTP servers or FTP drops."
  },
  {
    "technique_id": "T1056.001",
    "technique_name": "Keylogging",
    "family": "AgentTesla",
    "text": "AgentTesla installs a low-level keyboard hook using SetWindowsHookEx to capture all keystrokes on the infected system. The keylogger records input across all applications and periodically transmits logs to the command-and-control server. This enables theft of credentials typed into login forms and sensitive documents."
  },
  {
    "technique_id": "T1113",
    "technique_name": "Screen Capture",
    "family": "AgentTesla",
    "text": "AgentTesla periodically captures screenshots of the victim's desktop using GDI+ APIs. Screenshots are taken at configurable intervals and saved as JPEG images before exfiltration. This allows attackers to observe user activity and capture information displayed on screen."
  },
  {
    "technique_id": "T1071.001",
    "technique_name": "Web Protocols",
    "family": "AgentTesla",
    "text": "AgentTesla exfiltrates collected data over HTTP/HTTPS POST requests to web panels or Telegram bot APIs. It formats stolen credentials and keylog data into structured messages. The use of legitimate web protocols helps the traffic blend with normal browsing activity."
  },
  {
    "technique_id": "T1547.001",
    "technique_name": "Registry Run Keys / Startup Folder",
    "family": "AgentTesla",
    "text": "AgentTesla establishes persistence by adding entries to the HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run registry key. Some variants also copy themselves to the Startup folder. This ensures the keylogger and credential stealer restart after system reboots."
  },
  {
    "technique_id": "T1219",
    "technique_name": "Remote Access Software",
    "family": "Remcos",
    "text": "Remcos is a commercial remote access tool repurposed for malicious control of infected systems. It provides full remote desktop access, file management, and command execution capabilities. Operators use it to maintain persistent interactive access to compromised endpoints."
  },
  {
    "technique_id": "T1055",
    "technique_name": "Process Injection",
    "family": "Remcos",
    "text": "Remcos injects its payload into legitimate Windows processes such as explorer.exe or svchost.exe using process hollowing. The malware creates a suspended process, unmaps its memory, and writes its own code into the address space. This technique evades process-based detection and inherits the trust level of the host process."
  },
  {
    "technique_id": "T1056.001",
    "technique_name": "Keylogging",
    "family": "Remcos",
    "text": "Remcos includes a built-in keylogger that captures all keyboard input on the compromised system. Keystrokes are logged with window titles for context and stored locally before transmission. The operator can enable or disable keylogging remotely through the C2 channel."
  },
  {
    "technique_id": "T1571",
    "technique_name": "Non-Standard Port",
    "family": "Remcos",
    "text": "Remcos communicates with its command-and-control server over non-standard TCP ports, often using ports like 2404, 7832, or other uncommon values. The traffic is encrypted with RC4 to avoid content inspection. Using non-standard ports helps evade firewall rules configured for common services."
  },
  {
    "technique_id": "T1112",
    "technique_name": "Modify Registry",
    "family": "Remcos",
    "text": "Remcos modifies Windows registry keys to store its configuration data and establish persistence. Configuration values including C2 addresses and encryption keys are stored in encoded form under HKCU registry hives. The malware also disables security features by modifying registry-based policies."
  },
  {
    "technique_id": "T1059.001",
    "technique_name": "PowerShell",
    "family": "DCRat",
    "text": "DCRat leverages PowerShell for initial execution and post-exploitation activities. It downloads and executes additional modules using encoded PowerShell commands. The malware uses PowerShell's reflection capabilities to load .NET assemblies directly into memory without touching disk."
  },
  {
    "technique_id": "T1125",
    "technique_name": "Video Capture",
    "family": "DCRat",
    "text": "DCRat can activate the victim's webcam to capture video and still images. The captured media is stored temporarily before being uploaded to the attacker's server. This surveillance capability is part of DCRat's modular plugin architecture."
  },
  {
    "technique_id": "T1496",
    "technique_name": "Resource Hijacking",
    "family": "DCRat",
    "text": "DCRat includes a cryptocurrency mining plugin that hijacks CPU and GPU resources on infected systems. The miner runs as a low-priority process to avoid detection through performance degradation. Mining proceeds are directed to attacker-controlled wallet addresses."
  },
  {
    "technique_id": "T1055.012",
    "technique_name": "Process Hollowing",
    "family": "DCRat",
    "text": "DCRat uses process hollowing to inject its payload into legitimate Windows processes. It creates a suspended instance of a trusted executable, replaces its memory contents with malicious code, and resumes execution. This allows DCRat to operate under the guise of a legitimate process."
  },
  {
    "technique_id": "T1573.001",
    "technique_name": "Symmetric Cryptography",
    "family": "DCRat",
    "text": "DCRat encrypts its command-and-control communications using AES symmetric encryption. The encryption key is derived from a hardcoded seed embedded in the malware configuration. This prevents network monitoring tools from inspecting the content of C2 traffic."
  },
  {
    "technique_id": "T1071.001",
    "technique_name": "Web Protocols",
    "family": "Androm",
    "text": "Andromeda (Gamarue) communicates with its command-and-control infrastructure over HTTP using custom-encoded POST requests. The bot sends system fingerprint data and receives plugin commands in obfuscated binary format. C2 domains are hardcoded and rotated across campaigns to resist takedowns."
  },
  {
    "technique_id": "T1129",
    "technique_name": "Shared Modules",
    "family": "Androm",
    "text": "Andromeda uses a modular plugin architecture that downloads and loads additional DLL-based payloads at runtime. Plugins include keyloggers, browser form grabbers, rootkits, and credential stealers. This modularity allows operators to customize infections post-compromise without redeploying the core bot."
  },
  {
    "technique_id": "T1547.001",
    "technique_name": "Registry Run Keys / Startup Folder",
    "family": "Androm",
    "text": "Andromeda persists by adding registry run key entries under HKCU or HKLM and copying itself to the user's AppData or Temp directory under a randomized filename. Some variants also create scheduled tasks for redundant persistence. The bot re-establishes C2 communication immediately after reboot."
  },
  {
    "technique_id": "T1091",
    "technique_name": "Replication Through Removable Media",
    "family": "Androm",
    "text": "Andromeda spreads to removable USB drives by creating hidden copies of itself and placing autorun.inf files that trigger execution when the drive is inserted into another system. It hides original files on the drive and replaces them with malicious shortcuts. This worm-like propagation method enables offline spreading across air-gapped networks."
  },
  {
    "technique_id": "T1564.001",
    "technique_name": "Hidden Files and Directories",
    "family": "Androm",
    "text": "Andromeda sets hidden and system attributes on its dropped files to avoid casual user detection. On infected USB drives it hides legitimate user files and creates visible shortcuts that execute the malware when opened. The bot also injects into explorer.exe to maintain stealth during execution."
  },
  {
    "technique_id": "T1555.003",
    "technique_name": "Credentials from Web Browsers",
    "family": "SnakeLogger",
    "text": "Snake Keylogger extracts saved passwords, cookies, and autofill data from Chromium-based and Gecko-based browsers. It reads the SQLite login databases and decrypts stored credentials using the Windows Data Protection API (DPAPI). The harvested credentials are formatted and exfiltrated to attacker infrastructure via SMTP, FTP, or Telegram."
  },
  {
    "technique_id": "T1056.001",
    "technique_name": "Keylogging",
    "family": "SnakeLogger",
    "text": "Snake Keylogger installs a global keyboard hook to capture all user keystrokes in real-time. It logs active window titles alongside keystrokes to provide context for captured credentials. Keystroke data is buffered and periodically exfiltrated to the configured C2 channel."
  },
  {
    "technique_id": "T1113",
    "technique_name": "Screen Capture",
    "family": "SnakeLogger",
    "text": "Snake Keylogger captures screenshots at regular intervals to monitor user activity and capture visually displayed sensitive information. Screenshots are compressed and queued for exfiltration alongside keystroke logs. This provides attackers with visual confirmation of stolen data entry."
  },
  {
    "technique_id": "T1071.003",
    "technique_name": "Mail Protocols",
    "family": "SnakeLogger",
    "text": "Snake Keylogger exfiltrates stolen data via SMTP to attacker-controlled email accounts. Credentials, keystrokes, and screenshots are formatted into email messages with structured subject lines for automated parsing. Some variants alternatively use FTP upload or Telegram bot APIs for exfiltration."
  },
  {
    "technique_id": "T1027.004",
    "technique_name": "Compile After Delivery",
    "family": "SnakeLogger",
    "text": "Snake Keylogger delivery chains often use multi-stage .NET loaders that dynamically compile and execute payloads in memory. The initial dropper decrypts an embedded resource, loads it via reflection, and executes without writing the final payload to disk. This technique evades static file scanning and signature-based detection."
  },
  {
    "technique_id": "T1113",
    "technique_name": "Screen Capture",
    "family": "Mokes",
    "text": "Mokes captures screenshots of the victim's desktop at configurable intervals. Screenshots are saved as temporary image files before being queued for upload to the command-and-control server. This surveillance capability enables visual monitoring of user activity and confirmation of target value."
  },
  {
    "technique_id": "T1056.001",
    "technique_name": "Keylogging",
    "family": "Mokes",
    "text": "Mokes implements a keylogger that records all keyboard input on the infected system. Keystrokes are stored in local log files with timestamps and transmitted to the C2 server during periodic check-ins. The keylogger captures credentials and sensitive communications across all applications."
  },
  {
    "technique_id": "T1125",
    "technique_name": "Video Capture",
    "family": "Mokes",
    "text": "Mokes can capture audio and video from the system microphone and webcam. Media recordings are stored locally and exfiltrated to the C2 infrastructure on a configurable schedule. This gives operators real-time surveillance capabilities on compromised endpoints."
  },
  {
    "technique_id": "T1547.001",
    "technique_name": "Registry Run Keys / Startup Folder",
    "family": "Mokes",
    "text": "Mokes establishes persistence by placing copies of itself in the user's Application Data directory and registering them with startup registry keys. It uses randomly generated filenames to avoid simple signature detection. The backdoor re-initializes C2 communication on every system boot."
  },
  {
    "technique_id": "T1573.001",
    "technique_name": "Symmetric Cryptography",
    "family": "Mokes",
    "text": "Mokes encrypts all command-and-control communications using AES-256 encryption. The encryption keys are derived from hardcoded values embedded in the binary. This prevents network-level inspection of data exfiltration and command traffic between the backdoor and its operators."
  },
  {
    "technique_id": "T1036",
    "technique_name": "Masquerading",
    "family": "Benign",
    "text": "Legitimate software uses standard naming conventions and proper code signing certificates issued by recognized certificate authorities. Benign executables have consistent metadata including valid company names, version information, and appropriate file descriptions. The absence of name spoofing or metadata inconsistencies is characteristic of clean software."
  },
  {
    "technique_id": "T1007",
    "technique_name": "System Service Discovery",
    "family": "Benign",
    "text": "Legitimate system administration and monitoring tools query running services for health checks and status reporting. Clean software interacts with the Service Control Manager through documented APIs with appropriate privileges. These operations are bounded, logged, and performed with user consent or administrative authorization."
  },
  {
    "technique_id": "T1083",
    "technique_name": "File and Directory Discovery",
    "family": "Benign",
    "text": "Clean applications access file systems through standard APIs with user-initiated actions such as file open dialogs and project workspace scanning. Benign file enumeration is typically scoped to user-designated directories rather than searching for credential stores or sensitive data paths. Normal software does not systematically scan for browser databases or cryptocurrency wallets."
  },
  {
    "technique_id": "T1518",
    "technique_name": "Software Discovery",
    "family": "Benign",
    "text": "Legitimate software may query installed programs for compatibility checks, dependency verification, or update management. These queries are transparent, well-documented, and serve a clear functional purpose for the user. Unlike malware reconnaissance, benign software discovery does not attempt to identify or disable security products."
  },
  {
    "technique_id": "T1204",
    "technique_name": "User Execution",
    "family": "Benign",
    "text": "Benign software requires explicit user consent for installation and execution through standard operating system prompts. Clean applications display UAC dialogs, license agreements, and clear installation wizards. They do not employ social engineering, exploit execution, or bypass user interaction requirements to achieve code execution."
  }
]
